# hbase2.2/profile.d/java.sh
#!/usr/bin/env bash

export JAVA_HOME=/usr
export HBASE_MANAGES_ZK=false
